# profile/airootfs/etc/skel/.config/autostart/force-resolution.desktop
[Desktop Entry]
Type=Application
Exec=/usr/local/bin/force-resolution.sh
Hidden=false
NoDisplay=false
X-KDE-autostart-phase=1
Name=Force Resolution
